# docs-jekyll-backup/examples/atproto-pds.service
[Unit]
Description=ATProto Personal Data Server (PDS)
After=network.target

[Service]
Type=simple
User=pds
Group=pds
WorkingDirectory=/opt/nspds
ExecStart=/opt/nspds/bin/atprotopds-cli serve --data-dir /var/lib/nspds --config /etc/nspds/config.json --foreground
Restart=always
RestartSec=5

# Logging configuration
# The PDS handles its own file logging if configured in config.json.
# This ensures stdout/stderr are also captured by journald.
StandardOutput=append:/var/log/nspds/stdout.log
StandardError=append:/var/log/nspds/stderr.log

[Install]
WantedBy=multi-user.target
