# config/pds.service
[Unit]
Description=ATProto PDS Server (kaszlak)
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=${DEPLOY_DIR}
ExecStart=${DEPLOY_DIR}/build-linux/bin/kaszlak serve --port 2583 --hostname 0.0.0.0 --config ${DEPLOY_DIR}/config/production.json --foreground
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
