# ops/deploy/pds.service
[Unit]
Description=ATProto Personal Data Server (NSPds)
Documentation=https://github.com/jvalinsky/NSPds
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=atprotopds
Group=atprotopds

WorkingDirectory=/opt/atprotopds
ExecStart=/opt/atprotopds/bin/atprotopds-cli serve --config /etc/atprotopds/production.json

# Environment
Environment=PDS_LOG_LEVEL=info
Environment=PDS_LOG_FORMAT=json
Environment=PDS_USE_KEYCHAIN=true

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/atprotopds /var/log/atprotopds
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atprotopds

[Install]
WantedBy=multi-user.target
