# ATProtoPDS/Resources/LaunchAgents/com.atproto.pds.user.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.atproto.pds.user</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/september</string>
		<string>serve</string>
		<key>RunAtLoad</key>
		<true/>
		<key>KeepAlive</key>
		<false/>
	<key>ProcessType</key>
		<string>Background</string>
	<key>WorkingDirectory</key>
		<string>~/.config/september</string>
	<key>StandardOutPath</key>
		<string>~/Library/Logs/september/agent.log</string>
	<key>StandardErrorPath</key>
		<string>~/Library/Logs/september/agent.log</string>
</dict>
</plist>
